# .config/systemd/user/volatize-sync@.service
[Unit]
Description=

[Service]
Type=oneshot
ExecStart=%h/bin/volatize.sh start %i
