# .config/systemd/user/xcape.service
[Unit]
Description=xcape daemon

[Service]
Type=forking
ExecStart=/usr/bin/xcape

[Install]
WantedBy=window-manager.target
